# Crypto/Dad-Joke-Crypto/build/supervisord.conf
[supervisord]
nodaemon=true
user=root
[program:app]
command=python3 app.py
directory=/app
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log
