# Web/Pin-In-The-Bin/build/supervisord.conf
[supervisord]
nodaemon=true
user=ctf
[program:app]
command=/usr/local/bin/python3 app.py
directory=/home/ctf/app
autostart=true
autorestart=true
stderr_logfile=/var/log/app.err.log
stdout_logfile=/var/log/app.out.log
